# arxos-news.timer
[Unit]
Description=Check ARXOS news periodically

[Timer]
OnStartupSec=3min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
